# client/supervisord.conf
[supervisord]
nodaemon=false

[program:rsyslog]
command=rsyslogd -n

[program:cron]
command=cron -l 2 -f
user=root
