# bahub-client/supervisord.conf
[supervisord]
nodaemon=false

[program:rsyslog]
command=rsyslogd -n

[program:cron]
command=crond -d 6 -f
user=root
